# docker/php/crontab
PATH=/usr/sbin:/usr/bin:/sbin:/bin:/usr/local/bin
* * * * * cd /var/www/html && php artisan schedule:run >> /dev/null/ 2>&1